# server/com.kritsadas.spotify-cyd-mdns.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Advertises _spotify-cyd._tcp from the HOST so the board can discover
         the server via mDNS. Needed because the server runs in Docker, where
         multicast can't reach the LAN on macOS. Install:
           cp server/com.kritsadas.spotify-cyd-mdns.plist ~/Library/LaunchAgents/
           launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.kritsadas.spotify-cyd-mdns.plist -->
    <key>Label</key>
    <string>com.kritsadas.spotify-cyd-mdns</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/dns-sd</string>
        <string>-R</string>
        <string>server</string>
        <string>_spotify-cyd._tcp</string>
        <string>.</string>
        <string>8080</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
